# deploy/app-platform/supervisord.conf
# No unix_http_server — avoids /run/supervisor.sock stale-socket loops on fast App Platform restarts.
# supervisorctl is not needed in the container.

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:vendure]
directory=/app/vendure
command=/usr/local/bin/node dist/index.js
autostart=true
autorestart=true
priority=10
startretries=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
# App Platform sets PORT=8080 for nginx; Vendure reads PORT and must stay on 3000 (same as nginx upstream).
environment=NODE_ENV="production",PORT="3000"

[program:vendure-worker]
directory=/app/vendure
command=/usr/local/bin/node dist/worker.js
autostart=true
autorestart=true
priority=15
startsecs=8
startretries=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",PORT="3000"

[program:storefront]
directory=/app/storefront
command=/usr/local/bin/node ./node_modules/next/dist/bin/next start -p 3001 -H 127.0.0.1
autostart=true
autorestart=true
priority=20
startsecs=5
startretries=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=NODE_ENV="production"

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=100
startsecs=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
